# portfolio-dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aditya.portfolio.dev</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/adityakumar/my_Portfolio && npm run dev</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/adityakumar/portfolio-dev.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/adityakumar/portfolio-dev-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/adityakumar/my_Portfolio</string>
</dict>
</plist>

